# docker/app/supervisor.conf
[program:tron]
user=root
directory=/root/tron_backend/
command=uwsgi --http :8000 --module tron_backend.wsgi
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0